# dec2/formulas/salt/jupyter/templates/jupyter-notebook.conf
{%- from 'conda/settings.sls' import install_prefix with context -%}
{%- from 'jupyter/settings.sls' import jupyter_config_dir, notebooks_dir, user, home_dir, port with context -%}

{%- set environment = [] -%}
{%- do environment.append('PATH="' ~ install_prefix ~ '/bin:%(ENV_PATH)s"') -%}

{%- do environment.append('USER="' ~ user ~ '"') -%}
{%- do environment.append('HOME="' ~ home_dir ~ '"') -%}
{%- do environment.append('JUPYTER_CONFIG_DIR="' ~ jupyter_config_dir ~ '"') -%}

[program:jupyter-notebook]
command={{ install_prefix }}/bin/jupyter notebook --ip=0.0.0.0 --port={{ port }} --no-browser
user={{ user }}
numprocs=1
directory={{ notebooks_dir }}
autostart=false
autorestart=true
stdout_logfile=/var/log/jupyter-notebook.log
stdout_logfile_maxbytes=50MB
redirect_stderr=true
environment={{ environment | join(',') }}
