# dec2/formulas/salt/dask/distributed/templates/dworker.conf
{%- from 'conda/settings.sls' import install_prefix with context -%}
{%- from 'dask/distributed/settings.sls' import scheduler_host with context -%}

[program:dworker]
command={{ install_prefix }}/bin/python {{ install_prefix }}/bin/dworker {{ scheduler_host }}:8786
startsecs=1
numprocs=1
autostart=false
autorestart=true
stdout_logfile=/var/log/dworker.log
stdout_logfile_maxbytes=50MB
redirect_stderr=true
